# files/sdw-notify.service
[Unit]
Description=SecureDrop Updater GUI Notification

[Service]
Type=oneshot
ExecStart=/usr/bin/sdw-notify
